# example/example.socket
[Unit]
Description=example Socket

[Socket]
ListenStream=/var/run/example.socket

[Install]
WantedBy=sockets.target
